# fb-tw/bin/daysale.cron
#! /bin/sh
cd /data/website/island/bin/
 /usr/local/php-cgi/bin/php -f /data/website/island/bin/daysale.php